# deployment/deployment_files/InfluxDB_Grafana_backup.crontab
0 0 * * * /bin/bash /home/ubuntu/scripts/InfluxDB_Grafana_backup.sh
