# deployment/deployment_files/InfluxDB_Grafana_backup.crontab
0 * * * * sh /home/ubuntu/scripts/InfluxDB_Grafana_backup.sh
